# supervisord.conf
[supervisord]

nodaemon=true

[program:ife]

command=/usr/bin/python webserver.py